# contrib/systemd/pgwd.service
# systemd unit for pgwd (daemon mode).
# Copy to /etc/systemd/system/pgwd.service and adjust paths.
# Config: /etc/pgwd/pgwd.conf (or -config /path/to/custom.conf). No env file.
#
#   sudo cp contrib/systemd/pgwd.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pgwd

[Unit]
Description=Postgres Watch Dog (daemon)
Documentation=https://github.com/hrodrig/pgwd
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# .deb/.rpm install to /usr/bin. For manual install to /usr/local/bin, change path.
ExecStart=/usr/bin/pgwd
Restart=on-failure
RestartSec=30

# Optional: custom config path (default is /etc/pgwd/pgwd.conf)
# ExecStart=/usr/bin/pgwd -config /etc/pgwd/pgwd.conf

# Optional: run as a dedicated user (create with: useradd -r -s /bin/false pgwd)
# User=pgwd
# Group=pgwd

[Install]
WantedBy=multi-user.target
