# contrib/systemd/pgwd-once.service
# systemd unit for pgwd (one-shot, used by pgwd.timer).
# Copy to /etc/systemd/system/ together with pgwd.timer.
# Config: /etc/pgwd/pgwd.conf (or -config for custom path). No env file.
#
#   sudo cp contrib/systemd/pgwd-once.service contrib/systemd/pgwd.timer /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now pgwd.timer

[Unit]
Description=Postgres Watch Dog (one-shot)
Documentation=https://github.com/hrodrig/pgwd
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# .deb installs to /usr/bin; manual install often to /usr/local/bin.
ExecStart=/usr/bin/pgwd

# No Restart= for oneshot; the timer will run again on schedule.
