# contrib/freebsd/pkg-plist
bin/pgwd
share/man/man1/pgwd.1.gz
share/doc/pgwd/LICENSE
etc/pgwd/pgwd.conf.example
etc/rc.d/pgwd
